# config/bigdata/openresty1_19_3_1/script/taotao-cloud-access-log.cron
0 23 * * * sh /opt/openresty/script/spilt-access-log.sh >> /opt/openresty/script/logs/spilt-access-log.log 2>&1
